# install/crontab
*/1 * * * * php /var/www/html/php/index.php -c task -a run